# deployment/mcp-rag-server.service
[Unit]
Description=MCP RAG Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=rafalwolak
WorkingDirectory=/Users/rafalwolak/dev/mcp/rag
Environment=PATH=/Library/Frameworks/Python.framework/Versions/3.13/bin
Environment=GEMINI_API_KEY=AIzaSyBydSC6F2F2IEP_PWXsNljkKOqtOfrwj7g
Environment=QDRANT_URL=http://localhost:6333
Environment=MEM0_SELF_HOSTED=true
Environment=MEM0_LOCAL_STORAGE_PATH=./mem0_data
Environment=FASTMCP_HOST=127.0.0.1
Environment=FASTMCP_PORT=8001
ExecStart=/Library/Frameworks/Python.framework/Versions/3.13/bin/python run_server_http.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 